# register_vm_files/register_vm
#!/sbin/openrc-run
description="Register VM with Prometheus server"

command="/usr/local/bin/register_vm.sh"
command_background="no"
